# infrastructure/ansible/roles/genesis/templates/java-application.service.j2
[Unit]
Description={{ name }}
After=network-online.target

[Service]
ExecStart=/usr/bin/java -jar {{ jar_file }} {{ configuration_file }}
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
